# _addons/App/501/_mdl/501-image_reset.0.cron
#!/usr/bin/perl
package CRON::module;
use open ':utf8', ':std';
use if $] < 5.018, 'encoding','utf8';
use utf8;
use strict;

=head1 NAME

501-image_reset.0.cron

=cut

=head1 DESCRIPTION

Removes all not original files (also resets the application)

=head1 DEPENDS

=over

=item *

L<App::501::_init|app/"501/_init.pm">

=back

=cut

use App::501::_init;



sub execute
{
	my %env=@_;
	
	alarm(0);
	
	my $sql=qq{
		SELECT
			CONCAT(ID_format,'/',SUBSTR(ID,1,4),'/',name,'.',file_ext) AS file_path
		FROM
			`$App::501::db_name`.a501_image_file
		WHERE
			ID_format != $App::501::image_format_original_ID
	};
	my %sth0=TOM::Database::SQL::execute($sql,'quiet'=>1);
	while (my %db0_line=$sth0{'sth'}->fetchhash())
	{
		main::_log("$db0_line{'file_path'}");
		`rm $tom::P_media/a501/image/file/$db0_line{'file_path'}`;
	}
	
	return 1;
	
	my $sql=qq{
		SELECT
			CONCAT(ID_format,'/',SUBSTR(ID,1,4),'/',name,'.',file_ext) AS file_part_path
		FROM
			`$App::510::db_name`.a510_video_part_file_j
		WHERE
			ID_format != $App::510::video_format_original_ID
	};
	my %sth0=TOM::Database::SQL::execute($sql,'quiet'=>1);
	while (my %db0_line=$sth0{'sth'}->fetchhash())
	{
		main::_log("$db0_line{'file_part_path'}");
		`rm $tom::P_media/a501/image/file_j/$db0_line{'file_path'}`;
	}
	
	return 1
}



=head1 AUTHORS

Comsultia, Ltd. (open@comsultia.com)

=cut

1;
