# _mdl/110/110-weblog_email_searchengines.7days.cron
#!/usr/bin/perl
# áéíóú - USE UTF-8 !!!
package CRON::module;
use TOM;
use open ':utf8', ':std';
use encoding 'utf8';
use utf8;
use strict;
#use App::400;
#use App::400::SQL::a400;

#use Chart;
#use Chart::2D::lines;
#use Chart::2D::columns;
#use Chart::2D::ring;

#use Net::HTTP::LiteAgent;

use TOM::Net::HTTP::CGI;
use Int::charsets::encode;

use TOM::Net::HTTP;
use TOM::Net::HTTP::UserAgent;
use TOM::Net::HTTP::referer;

#use Utils::datetime;
#use App::1B0;
use TOM::Database::connect;

use strict;

sub execute
{
 my %env=@_;
 
=head1
	my $val='q=IT+kari%C3%A9ra&hl=sk&lr=&cr=countrySK&start=40&sa=N';
	
	my %FORM=TOM::Net::HTTP::CGI::GetQuery($val);
	
	print "$FORM{q}\n";
	
#	if (utf8::is_utf8($FORM{'q'}))
#	{
#		print "je utf8\n";
#	}
#	utf8::decode($FORM{'q'});
	if (utf8::is_utf8($FORM{'q'}))
	{
		print "SUPER je utf8\n";
	}
	
	$FORM{'q'}=Int::charsets::encode::UTF8_ASCII2($FORM{'q'});
	
	Int::charsets::encode::UTF8_ASCII2('čačané');
	
#	while ($FORM{'q'}=~s/^(.)//)
#	{
#		print "-$1\n";
#	}
	
	print "$FORM{q}\n";
	
	
 return 1;
=cut
 
 TOM::Database::connect::multi("stats");
 
 #print "connected to $TOM::DB{stats}{name}\n";
 my $boundary = "------------080308060200060803000107";

 my $date=`date -u "+%a,%e %b %Y %H:%M:%S (%Z)"`;chomp($date);

 $env{to_email} = $TOM::contact{'stats'}.";".$TOM::contact{'TOM'};
 #$env{to_email} = $TOM::contact{'stats'}.";".$TOM::contact{'TOM'};
 #$env{to_email} = "gregor\@webcom.sk;fordinal\@webcom.sk";
 #$env{to_email} = "gregor\@webcom.sk";

my $pics = <<"PICS";
Content-Type: image/gif;
 name="google.gif"
Content-Transfer-Encoding: base64
Content-ID: <part1.google.com\@webcom.sk>
Content-Disposition: inline;
 filename="google.gif"

R0lGODlhNwAUAPcAAIACBASDBKyGhNTGlAQChMTCzPTizARaBGSG1EROfMzi7MympCBIvOT0
9PRCJAQCZOnBBLTI8GwuFJSixPn77ERmuLyKjIzijJSGVAQkhNyilLhURFxqjOzq7NzWzAQj
pMTCrDxePIRydKcCBHSCrMz21ARG1ICs/OT7/BhCqLSqrLTW9OTCvISazNzy/LAoFLyCfExq
RLyGBOTm/AQkcDRy3AQpmdS+pOTY3MR2bLS2tGwSBN/l5Nnl/DRetNS6PKSyzPyWhGRypBQy
XKyelMzKzFRafPj99PzmpAQyrIySpGyS5FxuXMSmPIyNnPzy7IRiBL/P9EhwxPheSOzu7MzC
vMwCBPzy1BRU3KS25Pz+/MzW/OTSxMQyHLymXMSWBOyelAQalMDO0ARGzIQ+NGSuZAQuhOyi
lHx+jPzo7PzSzCxSpPH0/Px2XLy2vPzo3Bg0dPz+lJiepJx+JNDKwJwCBHRWBFSC5CROlNzk
8PzeJHQqJIyi3HSKdIxedOTa1AQuxPzuPCRCbHxeZLzyvPzaZNTCNJR+XJxyBMy+lKyCdOza
jPSGdLxiXNS+zFx+1JQaBLRuZBxaHMyurPS6rNwyJNSalOzivMSudExejJwmFGR8YFyS9IyO
jHwaDAQcqPz6rAw+tKR+PKSSdPz+1GRWPBxa5OR+bCQ+ZNyopNTa7PyqpERWjOTy5ISatFRq
VNy2tBxStJyuzLQOBMiQBAQylMwOBNSgBBRSzKTG/KimvDRaxKS23LjM3HF7nKyShCxKpPzK
vHR6rCRa1AQupJx+BEx+5AQGlHSGvIyWtIye1BQqbHia6LwyLMyuZBA2gPzqtAQyvPzy5Pz+
5PzGuMy5sfzOvDxuvODe5OQ6HIyy/Ofu/KS+5PzibPSOfFiE3Lzm9OTU1LS21O+XjLTO/MTE
xOx6ZLy6nKRqbHwKBExWfAxO1Njq/FR2XDRapJSq3ASKBBxiHPx+bFRyVNzK1DxelNzKnPLl
1Ojo9HyW1Pn19ExyzGwuLPy+tHQeFMze9CH5BAAAAAAALAAAAAA3ABQARwj/ALUIHEiwoMGD
CBMq/DcGzhGFbXYQDDJCFylaTghCQMSmRoKBpuAIDNBOizl+FOTlECBHocuXAo09ezgwhQqY
AxkYeVOQjRMPOIPCiwETiI0P+QpGGyJQUx1YWoJ1GSTQHqImWo4UmqNr4BFmtXgJLHEg3sNm
IyoFGbeATUEXCHwIcyIkFCAEeYLCZANPksttnJalcemt0aQFL2yxEDjNUCIQmG79GLjk3YR1
KZINDjBP4CkYOiJZAdOgAxUtxj4gI6ivhat1zxpkA4BioD4rIkDRGkUwkJ1LzAjI0sIOWKeR
TLRMGaHgpbRyVQq8eUJBH0E2bBocnIFvm0A2FYTQ/7lHx10YerYHsuHhXa/7gTo44FkjrMqV
9wNXaOPWHn9CNoRcUAZM2+ziCxdPHEHBGyQ8oIR/CcETwku50HDDS2poMM4qNA2ExCLdQFMQ
PltsE0UWM2gBzysDUTMOGNRogU0RVPSAhSAuUQNJLwKhMEIdAkFwyBYCDUCLF1qc4I4B0/Ty
iS8jddbPC+AINIsneWmRxzAZ4HBdLlv0EIEWDnhSm0BvjOBHHDIQQZAeUKixjhkDSYHKSCxO
kU41YmTJWhIfZGGQCTgyMoIlA1ECgBha3FKMW0GW4p0crNzBQDUUjLSJFqmMsIFC2GSSgQ2x
7CNFBeqgIZALipCxgQVuqEQyED0YiOLMOeJ8Rw4fFRDziTJatNLHJo5o4cIve+yBTjkQ4qQN
DQYQRExGzVbLyzG7kGPiI76EU+23WnSADQ88WBdUQAA7
--$boundary
Content-Type: image/gif;
 name="zoznam.gif"
Content-Transfer-Encoding: base64
Content-ID: <part2.zoznam.sk\@webcom.sk>
Content-Disposition: inline;
 filename="zoznam.gif"

R0lGODlhSAAUAPcAAMwCBOSCfNxCRPTCvNQiHOyinPTi3NxiXOyytMwSDNQyLPTS1Pzy7ORy
bNxSTOySjOy6vNQWHPzq5NxqZNQ2PPTa1OSKjNQmLOymrNxaXMwKBNxKRPTKxPzm7PTe5PTi
5ORmbOy2vPz69OR6fOyalMweHNQ+PNxiZMwSFPS+xPTa3OyurOSChPTGzNQqJOy2tOR2fNxW
XOy+vNxqbOSOjNQuLNxeXMwOFNxOVPTm5MwWFPTe3MwGDNxGTPTCxNQiJOyipNQyNPzy9ORy
dNxSVOySlMwaHPzq7NQ6PMwKDNxKTPTKzPz6/OyanNQeJNw+RORmZOyutOSGhORubMwGBNxG
RPzi3ORiXPSytPTW1PS6vOyKjNQqLOyqrORaXPze5Pzi5PS2vOR+fNQeHNw+PORiZNQSFPza
3PSurPS2tPS+vORqbOyOjNQuNOReXPzm5NQWFPze3PTGxNQmJOympNQ2NPz29OR2dNxWVOyW
lNQaHPzu7Nw6PMwODNxOTPTOzPz+/OyenOyGhP8AACD4AOYDABIAAAAAAGcUBzoAAG4AAHkA
AAAAAAAAABMAAAAAAAAAAAAAIAAAFAAAADiZAJlaABWMAACBABYAAAAAAAAAAMAAACBnAFwU
ABPtAAC/AAA0AADnAAASAAAAAH5n8AAUKQBubcB5eQAA/wAA/wAA/wAA//8ATP8A6P8AEv8A
AP8gT/9cAP8TAP8AAADhTwAbAABtAAB5AACiqAAU7QBuEgB5AABU/wDn/xMS/wAA/2vgAAAb
AABtAAB5AAAAXDQA6G4AEnkAACCwHlwiqxMUSwAAAGPSPDQT6G5uEnl5ABb/ewD/wAD/TsD/
AAFYcOfo6BISEgAAAAD70gETEwBubgB5ecCkK+XK6RJOEgAAALAMr//i/xJO/wAAf7Dg3P8b
6BJtEgB5AFVUsB/nIvgSFHcAACgA0iUAE/gAbncAef8AsP+/Iv9PFP8AAAAAAecwABIAAABl
AF9hABtyAG5jAHloAABfzgBlcBNuRwBnACH5BAAAAAAALAAAAABIABQARwj/AAEIHEiwoMGD
CBMqXMiQYIQuQgBZwZOEip8PGIwA4AJISxIAArR0kGLkDoMhSWr88VECBZ5AaZYwYUNlCgMx
SQTk6KLnhxo5cwT2IbIGCRwEDbgETULDDqCnUKHuifMHQgE2d2z4+UHFjA4zfXhQobKQioYk
fSKMuRBEAJ4pUoDIWHCEiYgPXBoOpIJEShEaUlgIHvBUhBwkZ3nspVJRLAAeipMoFsi48kDF
jDHz+Kh3YJIpOQB9GNHH4I0hCyT8aTJGz58dYHaoKFIjhSCBU5bMgRMIjIoXcoDomYMlQBIz
NOjo4WFDTh2ffFYgWIHFjxwlcqaYZsPAjh0RdoQI/zkSZ8WIKnrIdl7PvmGSIXGYCCnyQ8GL
Lx+ODBAQocgOD3swkAcSC3yQxQd7bNHHGCt8AdsHQwAwxwexveFBFAvksEcRcBzUB2cIkfEH
E3akMYd6SYywxwhJoGDBHkNM8EYISHx0wwt7RAAAHkcIYsIXXXDRhxRg+GHECgsEkaIEblCh
xAdFlAYAGyo0sUYbf7yBxQ8DofCAEEz8IUBBVGzm1RhGRECWGXO0UYcAPTiQwQlrTDHEHWKI
MYKeYtxxxxBTzHCCDTH4IQAfbfxwg3oM6VCEHHK80MQdRASBAqPtZYoQFSXUwN4cMZQxAwgn
ECEAFzdoqqqmfRTBRFSwMv/xBmmr1tpeEn4sAIgdXQTxgxd+BlEmEnYOMQQSfeCBBxEjTKBR
EkrY2QAROjx2gB8njOCHDlWMMEUJtRJggVMqNNCHHgi8EYAJFjARSBJzVDHFDkdkYMMRS9hw
wg4tzMFDDQ6UwQIgQACghB0L4HHHHm/cgYcKS3haUB8WgKAQD1UQxsQKFAiU4h5TkPUEwUY0
cMQfGTCWhxB1gGSAcAUcgcUWUTBwwg15GKAEFSd0gFMQC/AkEAt/ABEIEmoUAYENpk0QByB7
3FHtXhvs8EYgEOyRBxcfAJJGE0VAMUcFSwh0BxN4KJEDB1MUkcMCdbjQghZj3EDCFyZQAYIE
MyTDUcdQdyCQBR5ZvJBFDwQZIcUeYVYBYkFJGPGDEaVRYcQYJZQwBgoajKFDEklEQMANSeih
1OVGnKWHlCh0CEDpjCFxwxx4+LF5EH1wwSUAerz4lBplnNAAC1s00QUEavigvA9Y2JAEC3R0
sQIaK3RBx/V0AAHE9V2ggUAaWqQwQAsc/PFHFlmokMUCfwwAhBhKSJwQDxkc8ZRdQuwBhhyT
EmECF54jHaZsdZk+mKE+SriDFsAwg4uNYQ4/2NwACUjBhAQEADs=
--$boundary
Content-Type: image/gif;
 name="msn.gif"
Content-Transfer-Encoding: base64
Content-ID: <part3.msn.com\@webcom.sk>
Content-Disposition: inline;
 filename="msn.gif"

R0lGODlhMAAUAPcAAAQWTBSOLAQiZCRabHSOpLTK1AY+bAxelJyurARChCRqXOTm5ESOxDR2
lFR6bBROdPwuBFTCRISSlKyutBR2ZPz67JSipEhqhGR6lDxCZPyaJAQ2hAxGdAROhOTa1AQy
ZPxOBKy+xBRahPTu5DxelPRSHISWrBQ2XCxOdMQ+LFyGnPx+FHSyzHyavFxyfPw6DAdGfARS
lAQydLS+zDxihBSOTBxqpCx2vCROfIyOVMTS1Py+DIyitChWhNje5Jx6fAQqXARAfAxmrARK
nOTo9BRWfPz+/CxmhOju9PxqFHyetGSKpFx+jCRelARGhER2lAtOfPQuFISapCxqjPyuHBlK
bARPjMDGzPzWBAQeVCSmNAQuZEyepAQ/dJSuvOjm7GyevLS2vPz+9KSqrAQ6bCxdhGRGVHSI
lJSuxJyqvAQ6dMjW5Nzi5CxefFSStPw2BHyUpLSwvPn29JyirExyfHSClBROhAQ0bMS+vOzu
7DxunHyWtLR6VHyazFx2jPA6FAVKfEyyTHyWXCyKhNSOTCSCRJx2ZHyqZOyCRHlGTOyqLGSo
1PQ+HEyyPHzORARifPzODPxOFNSeTERypCxOfCRytPRmJEx+pCSClPyKHPzeBERWbGR2hExO
XDR/vFxulPSeNAQ6lJCivNTi3DRmlNTGxAQudHCStERidBRilBSWNPzGFIyqtIyerPy2FJ62
xMTK1HRufJyepFSSzLTC1Nza7LTGxBxWfCyuNAQyXPQ2FPRuJBRanKSytHyWnGR+nARKjER2
nHyGnPRKHAQ+jPRSJAQ2fMTCzJymtOTa5GyCjCRinBROfDxujPzSDAQfXKSwvPTq7CxejKSy
xPzy7Nza3BRKdLzG1JSetAQWVAQibBRajARSnCxmjPxqHHSKpAQubAQ6fNTW5Pw2DARCbERu
hGx+lERGZAlKdAxShOzy5ICarGSKnEBmhMzW1Ozy9GSCjAVKhEx6lAxSfOzq7Pz69KCmrPQ+
FMzO1PwyBKyytPyeJAQ6hAQ2ZPRWHAQ2dCH5BAAAAAAALAAAAAAwABQARwj/ADsI7ABFYEGC
BwcaTKjQ4MCEBwuma9ghXrwuUmT54oDM3SsYRdKsCYHCFAolyYgQ8bLtH445FoR1wbYmDapF
fWooIEEIlCV/IIYM5PDNiFFydYwagfJLqbVXSpVe2WP0HboLSrsdYxWvw4MOVig2DGvlQViB
Yb8+HKj2bAd00jxV2pbOCbBN3iLp+vMm0cB4/fAYvSIAqpECAvQZ1YE1aj0vMq4YDZHNntFR
4RK0CYSLgiFFkhrwIwjWrFe3pwmqVevQoFt0FQV+BRaj0598KWLILvgVduyuwB50lS0wHmvW
8dChY7KIwQ0hvBBpsDTsRT6/Aw0kNUovah6l9OAY/11TRZhRNOKNfonaw0EjLYN+uKKyCgsx
0jBOKMUApILRJeCsZwQB4ZRBQzunIGGUPt8Z0YQMAtrzTzzAwICDFYA4AYgKMPwVhDH/BCFi
ODKIGASJIiZgwAe5tGiMiCAG4cSJ/4QTBAxODCfWjqadxRpZx7HVUFtfWRHDEDGkY8VBq3lF
nJNgkcakj04WBAUvgADSlUUZQJDPC30xxOSTrCVkWmtqxTOPG57cYMMBt1jiDSPj6NJXQ07A
MGIXXaghohp9BhHPh3fcQUaJg5oYRBf/GKAEC7Ngsk0sGmSySwkQZMBWPKQopaBRpYghVQf1
RGXUFPJ4qtQ1n3ChSgAP8MixQjGR7BCKQjJIYRQR/WBjVAUfmGBUHtNUxQMGZfTzjGFsAOCL
UXKIEMQhjkQgCCRYYJHDBgMBYgAbRrXwTDJGiQIEPka9ksUe6pjajA9GXZJFAUbRQs6MAxTy
SAJOTCKOHgMlgI5S2xig1BFllGrEGaIigUE54BohilFizFOFUsEEAYUTHQDzFwddFYQOJ9Ag
wIMaLsRhchd0hDGGF2RUYYIH9dTjwSkPMDHBBEv8QwMC0ECDTgPdQGHWA18xw9qOTDftNNMB
AQA7
--$boundary
Content-Type: image/gif;
 name="webcom_logo.gif"
Content-ID: <part4.webcom.logo\@webcom.sk>
Content-Transfer-Encoding: base64
Content-Disposition: inline;
 filename="webcom_logo.gif"

R0lGODlhlgAcAMQAAL/Awd/g4O0bI/f395iam7i5urCxss/Q0YiKjJCSk+/v8NfY2Ofn6Kip
q8fIyfaNke84P/zU1v3i4/RxdvebnvV/hKCho4CChP///wAAAAAAAAAAAAAAAAAAAAAAAAAA
ACH5BAAAAAAALAAAAACWABwAAAX/ICaOZGmeaIpaV6u+cCzPL0MQFuLQfD+yLp9wSBQxEg6A
I0DYFZ8l4AVKrZ4aAIChYGEgrFApeEy1BACWRgGTIBPF7rjPrDUYMF95D67vwwBYAEcOFn40
fIaJJQM3Cwo5DIoxiJKSA4QWBQqVMJScnxgWohabJwOiDSmPowMkCgAELS0EmidSCgYJshcE
AK0qCgWxuxYOvyXBurIJtSdcoiIHDcsGpdENCC0J1SrTLWsnB7ILKAWyqSMO2buyOiZS7LsI
BygD5vGyhST2+LIFxz9k2cA3D8OjfvRQBGgHcMSwC/oWrbsQYISBfrvABcS4K+G+h/0iigDZ
j8C7dhwP/0zERw7FQ48jFu6yNkLcrHS7EgAIsKDAypYi4EHMYmAlxY/sEoxSJhKDzRaBvLWA
GZSdgQM8L8ZrgHWBVm0pnrYp8bXFnRLKLiRUMNEAQEbLSEjhNoJf3BH3ZlUkwcBAU6cQaR6k
WrUFgr042wEV4WBXJBQTFyvAB3DBXQwALrvaBdSTiMzjREyWRaBhDRERHjyIIAJxFIEnIp9I
67oEaIgkbksBIFcW75GyCIuQqtGziLTg8iIwHWOCgOcCJqgwLqV24Qu/Tc0coWz5urEYGDAc
sctAlvPnpUQ0jiGvPikaaTyADv1BCuqyrGP4mv3E17NP8ZZXQv+RwBE+68mSgv9lYOGRnw8Q
0PccBPcpaMuDJtzWnwmjtdDKLQblg8EA25F3IDsJBnGCTEHsoh8MEkJXoYqvtaCfhi/sJh5U
G13AwG3omDhVAEQWaSSRj4Vi4QkMXjDWRC++EKGEFK6wZI1HZejbCyxugyEGTRb1pQiyzcCe
e9fFN8N8EtpnJY292XgCji+kpRkb8YA3glRnycAeciIoxxwMFdBXwXRXxpmlbVu+8JQsToxw
W3AmPCqcCPZkZ9yke3XIy6AoLNCKBBRQIMGIiyn65qK5NarCACsttwg7eZhgZ6Qj9JWNprI0
ZNdNFu1CgHULBFkTL64xodaFcPZ4o6sqTNonCWVtOEKGk4YFksMuvJJWQBbYsJPkiCQppQUl
jxJgB0iE4Senli1Ye4Kn48a0C6iPctTtiYTBdWAJJOFjEparPhuvDFr95ZBZLyxgVFIwCSVw
vXXl1Q+8/ah5XcFzQqvCjpcCdgFN9RCSlAGpitCTxOpGOYICJhNTQMoGOQASLSTjNUo5o1Ds
1CjChQAAOw==
PICS

	$env{lng} = 'sk' unless $env{lng};
	my %TEXTS = (
		'sk' => {
			#header_subject => "[STAT][a110/weblog] last 7days domains/search engines",
			header_subject => "CYCLONE STATS: Search engines / 7 days",

			main_header => "Štatistiky návštevnosti subdomén <%domain-name%>",
			term => "za obdobie ",
			info => "Štatistiky Search engines nám ponúkajú presné informácie o tom, odkiaľ a ako užívatelia navštívili váš web. K dispozícii je zoznam domén a vyhľadávačov, odkiaľ návštevníci na web prišli, ako aj počet prístupov (reálne i pomerné číslo) z danej domény, resp. vyhľadávača v uplynulom období. Súčasťou je aj výpis všetkých kľúčových slov a fráz a počet vstupov na stránku prostredníctvom uvedených slov a fráz.",
		},

		'en' => {
			
		}
	);

	my $lastday = $main::time_current - 86400;
	my %last_date=Utils::datetime::ctodatetime($lastday,format=>1);
	$lastday = "$last_date{year}-$last_date{mom}-$last_date{mday}";
	my $last_datetime_string="$last_date{mday}. $last_date{mom}. $last_date{year}";
	my $today = $lastday; $today =~ s/^((\d+)-(\d+)-(\d+))$/$4. $3. $2/;

#=head1
my $email=<<"HEADER";
From: "$tom::Hm ($TOM::hostname)" <TOM\@$TOM::hostname>
To: <%TO%>
Subject: <%SUBJECT%>
Date: $date
List-Id: TOM3
MIME-Version: 1.0
Content-Type: multipart/related;
 boundary="$boundary"

This is a multi-part message in MIME format.
--$boundary
Content-Type: text/html;charset="UTF-8"
Content-Transfer-Encoding: 8bit

<html>
  <head>
    <meta http-equiv="content-type" content="text/html; charset=UTF-8" />
  </head>
  <body>
		<style>
			body { font-family: Arial, Verdana; }
		
			h1 {
				color: #808285;
				font-family: Arial, Verdana;
			}
		
			.info-domain { font-size: 18px; }
			.term { font-size: 14px; font-weight: normal; }
		
			h2, th { font-size: 12px; }
			th { text-align: left; }
			td { font-size: 10px; }

			#page { width: 500px; }
			#domains, #search_engines, #keywords, #keyphrases { width: 250px; float: left; }

			.odd { background: #FFF6F6; }
			.odd .first { border-left: 5px solid #D1D1D1; padding-left: 10px; }
			.first { border-left: 5px solid #fff; padding-left: 10px; }
			.endblock { clear: both;  }
		</style>

		<div id="page">

			<h1>
				<img src="cid:part4.webcom.logo\@webcom.sk" alt="webcom logo" border="0" /><br /><br />
				<span class="info-domain"><%main-header%></span><br />
				<span class="term"><%term%></span><br />
				<span class="term"><%info%></span>
			</h1>
			
			<%BODY%>

		</div>
	</body>
</html>
--$boundary
$pics
--$boundary--
HEADER
#=cut

 my $doc;
 my $doc1;
 
 
 # kolko % z uvodnych stranok tvoria prichody z vyhladavacov
 # kolko % zo stranok tvoria prichody z vyhladavacov
 #return 1;
 
 my $analyze=0;
 my %domains;
 my %domains_search;
 my $domains_search_count;
 my %keys;
 my %keys0;
 my $from=time-(7*86400);
 my %from_hash = Utils::datetime::ctodatetime($from, format=>1);
 my $db0=$main::DB{stats}->Query("
 	SELECT domain,domain_sub,referer
	FROM TOM.a110_weblog_rqs
	WHERE domain='$tom::Hm'
			AND reqtype='B'
			AND reqtime>$from
	ORDER BY reqdatetime ASC");
 #return 1;
 while (my %db0_line=$db0->fetchhash)
 {
  #return 1;
  #next unless $db0_line{referer};
  #next if $db0_line{referer}=~/$db0_line{domain}/;
  
  my ($domain,$query)=TOM::Net::HTTP::domain_clear($db0_line{referer});
  #$query="?".$query;
  
  $domains{$domain}++;
  
  if (my $dom=TOM::Net::HTTP::referer::analyze($domain))
  {
    #print "-i know domain $dom + $query\n" if $analyze;
    if (
         ($TOM::Net::HTTP::referer::table{$dom}{domain_type} eq "search engine")
	 &&($TOM::Net::HTTP::referer::table{$dom}{keywords_param})
       )
    {
   
     $domains_search{$dom}++;
     $domains_search_count++;
    
     my $keyword_param=$TOM::Net::HTTP::referer::table{$dom}{keywords_param};
     #my $regexp=$TOM::Net::HTTP::referer::table{$dom}{regexp_keywords};
     #my $keywords=($query=~/$regexp/)[0];
#     print " -parsing $dom-$query\n" if $analyze;
#     print " -$keywords\n";
		
		main::_log("query='$query'");
		my %FORM=TOM::Net::HTTP::CGI::GetQuery($query);
		next if $FORM{$keyword_param}=~/^cache/;
		next unless $FORM{$keyword_param};  
		
     $FORM{$keyword_param}=Int::charsets::encode::UTF8_ASCII($FORM{$keyword_param});
     $FORM{$keyword_param}=~tr/A-Z/a-z/;
     $FORM{$keyword_param}=~s|"(.*?) (.*?)"|"$1<M>$2"|g;
     $FORM{$keyword_param}=~s|["+&]||g;
     $FORM{$keyword_param}=~s| |;|g;
     $FORM{$keyword_param}=~s|^;||;$FORM{$keyword_param}=~s|;$||;$FORM{$keyword_param}=~s|;;|;|;
     $FORM{$keyword_param}=~s|<M>| |g;
     
     main::_log("keywords='$FORM{$keyword_param}'");
     
     next unless $FORM{$keyword_param};
     
     #$keys0{$FORM{$keyword_param}}++;
     
     my @keywords=split(';',$FORM{$keyword_param});
     
     my $par;
     #print " -[$dom]parsed:";
     #print " -parsed \"@keywords\"\n";
     foreach (sort @keywords)
     {
      #print "$_;";
      $keys{$_}++;
      $par.="$_ ";
     }
     #print "\n";
     $keys0{$par}++;
     
    }
  }
  
 }

 $email =~ s|<%main-header%>|$TEXTS{$env{lng}}{main_header}|g;
 $email =~ s|<%domain-name%>|$tom::Hm|g;
 $email =~ s|<%term%>|$TEXTS{$env{lng}}{term}$from_hash{mday}. $from_hash{mom}. $from_hash{year} - $today|g;
	$email =~ s|<%info%>|$TEXTS{$env{lng}}{info}|;

 $doc.="\n";
 
 my $count;
 #$doc.="domains (TOP 50):\n";
 $doc .= "<div id=\"domains\"><h2>domeny</h2><table>\n<tr><th>nazov</th><th>navstev</th></tr>\n";
 foreach (sort {$domains{$b} <=> $domains{$a}} keys %domains)
 {
  $count++;
  last if $count>50;
  #next if $domains{$_}<2;
  #$doc.=" +[$count] $_ $domains{$_}";
  my $subdoc = "<tr<%class%>><td class=\"first\"><%domain%></td><td>$domains{$_}</td></tr>\n";
  
  if (my $dom=TOM::Net::HTTP::referer::analyze($_))
  {
  	#$doc .=" - $dom/".$TOM::Net::HTTP::referer::table{$dom}{domain_type};
  	$subdoc =~ s|<%domain%>|<abbr title="$dom/$TOM::Net::HTTP::referer::table{$dom}{domain_type}">$_</abbr>|g;
  }

  else
  {
  	$subdoc =~ s|<%domain%>|$_|g;
  }
  $doc .= $subdoc;
  if ($count%2) { $doc =~ s|<%class%>| class="odd"|; }
  else { $doc =~ s|<%class%>||; }
  #$doc.="\n";
  
  
 }
$doc.="</table></div>\n";
 
 
 
 $doc.="\n";
 
 my $count;
 #$doc.="search engines [$domains_search_count] (TOP 50):\n";
  $doc .= "<div id=\"search_engines\"><h2>vyhladavace</h2><table>\n<tr><th>nazov</th><th>navstev</th></tr>\n";
 foreach (sort {$domains_search{$b} <=> $domains_search{$a}} keys %domains_search)
 {
  $count++;
  last if $count>50;
#  next if $domains_search{$_}<2;
  #$doc.=" +[$count] $_ $domains_search{$_} ".((int(($domains_search{$_}/($domains_search_count/100))*100))/100)."%\n";
  my $pict;
  if ($_ eq "google.com") { $pict = "part1.google.com\@webcom.sk"; }
  if ($_ eq "zoznam.sk") { $pict = "part2.zoznam.sk\@webcom.sk"; }
  if ($_ eq "msn.com") { $pict = "part3.msn.com\@webcom.sk"; }
  $doc .= "<tr<%class%>><td class=\"image first\"><img src=\"cid:$pict\" /></td><td>$_</td><td>$domains_search{$_}</td><td> ".((int(($domains_search{$_}/($domains_search_count/100))*100))/100)."%</td></tr>";
  if ($count%2) { $doc =~ s|<%class%>| class="odd"|; }
  else { $doc =~ s|<%class%>||; }
 }
 $doc .= "</table></div>\n";
 
 
  
 $doc.="\n<hr class=\"endblock\" />\n";
 
 my $count;
 #$doc.="keywords (TOP 50):\n";
  $doc .= "<div id=\"keywords\"><h2>klucove slova</h2><table>\n<tr><th>nazov</th><th>navstev</th></tr>\n";
 foreach (sort {$keys{$b} <=> $keys{$a}} keys %keys)
 {
  $count++;
  last if $count>50;
  #next if $keys{$_}<2;
  #$doc.=" +[$count] $_ $keys{$_}\n";
  $doc .= "<tr<%class%>><td class=\"first\">$_</td><td>$keys{$_}</td></tr>\n";
  if ($count%2) { $doc =~ s|<%class%>| class="odd"|; }
  else { $doc =~ s|<%class%>||; }
 }
 $doc .= "</table></div>\n";
 
 
 $doc.="\n";
 
 my $count;
 #$doc.="spojenia (TOP 50):\n";
   $doc .= "<div id=\"keyphrases\"><h2>frazy</h2><table>\n<tr><th>nazov</th><th>navstev</th></tr>\n";
 foreach (sort {$keys0{$b} <=> $keys0{$a}} keys %keys0)
 {
  $count++;
  last if $count>50;
  #next if $keys0{$_}<2;
  #$doc.=" +[$count] $_ $keys0{$_}\n";
  $doc .= "<tr<%class%>><td class=\"first\">$_</td><td>$keys0{$_}</td></tr>\n";
  if ($count%2) { $doc =~ s|<%class%>| class="odd"|; }
  else { $doc =~ s|<%class%>||; }
 }
 $doc .= "</table></div>\n";

 
 $doc.="\n";
 $email=~s|<%BODY%>|$doc|;
 #$email=~s|<%FROM%>|"$CRON::core_uname_n($tom::H)" <CRON\@$tom::H>|;
 $email=~s|<%SUBJECT%>|$TEXTS{$env{lng}}{header_subject}|;
 use Utils::datetime;
 #$email=~s|<%DATE%>|$Utils::datetime::DAYS{en}[$cron::Twday], $cron::Tmday $Utils::datetime::MONTHS{en}[$cron::Tmom-1] $cron::Fyear $cron::Fhour:$cron::Fmin:$cron::Fsec +-200|g;

 
 my %env0;
 foreach (split(';',$env{to_email})){$env0{$_}++;}
 $env{to_email}="";foreach (sort keys %env0){$env{to_email}.=$_.";" if $_}$env{to_email}=~s|;$||;
 $env{to_email_parse}=$env{to_email};$env{to_email_parse}=~s|;|>,<|g;$env{to_email_parse}="<".$env{to_email_parse}.">";
 $email=~s|<%TO%>|$env{to_email_parse}|g;


 #print "$email\n";
 
# print "ááááááááááá";

#=head1
#return 1;
	if
	(
		# sendtime='$main::time_current'
		# priority=0
		# from_name='CRON'
		# from_email='TOM\@$TOM::hostname'
		# from_host='$tom::H'
		# from_service='a110' a400 a300
		# to_name = ''
 
		$main::DB{main}->Query("
			INSERT INTO TOM.a130_send
			(
				sendtime,
				priority,
				from_name,
				from_email,
				from_host,
				from_service,
				to_name,
				to_email,
				body)
			VALUES	(
				'$main::time_current',
				'0',
				'CRON',
				'TOM\@$TOM::hostname',
				'$tom::H',
				'a110',
				'',
				'$env{to_email}',
				'$email'
				)"
		)
	)
	{
		main::_log(" ok, email sended");
	}
	else
	{
		main::_log(" :((, email not sended");
	}
#=cut
 return 1;
 }



1;
