# _mdl/130/130-control.forms-test.cron
#!/bin/perl
# ÁÉÍÓÚ - USE UTF-8 !!!
package CRON::module;
use Net::HTTP::LiteAgent;
use Net::SMTP;
use CVML;
use strict;
use Mysql;
use Digest::MD5  qw(md5 md5_hex md5_base64);
use MIME::Base64;

our $authors = "gregor\@webcom.sk";

# tento cron kontroluje emaily ulozene v databaze, pochadzajuce z adresy forms@webcom.sk
sub execute
{
	my %env=@_;
	main::_log("130 FORMS CONTROL");

	# DESIGN
	my %XSGN;
	
	$XSGN{mail} = <<"XSGN";
From: a130-forms <dev\@webcom.sk>
To: <\%email%>
Subject: [FORMS] ZLYHANIE APLIKACIE a130 ( <\%domain%> )

Toto je automaticky generovany email aplikacie na kontrolu formularov.
Z webu bol odoslany email, no nebol doruceny na forms\@webcom.sk
Je mozne, ze prislo ku problemom pri odosielani.

Obratte sa preto prosim na oddedelnie Developmentu.

Identifikacne udaje
-----------
ID-HASH: <\%hash%>
domena: <\%domain%>
subdomena: <\%domain_sub%>
sluzba: <\%from_service%>

Ine info
-----------
Cas odoslania: <\%sendtime%>
IP adresa odosielatela: <\%IP%>
mail bol odoslany na: <\%to_email%>
XSGN

	# oznacim dorucene emaily
	my $sql_mark = "
	UPDATE
		a130_send_history LEFT JOIN
		a130_received ON
		a130_received.body LIKE concat('".'%'."VERIFICATION-HASH: ',a130_send_history.hash,'%')
	SET
		a130_send_history.active='Y',
		a130_send_history.active='Y'
	WHERE
		a130_send_history.hash!='' AND
		a130_send_history.active='N' AND
		a130_received.active='N'
	";
	$main::DB{main}->Query( $sql_mark );

	# vyselektujem si nedorucene emaily
	my $sql = "
	SELECT
		a130_send_history.hash,
		a130_send_history.domain,
		a130_send_history.domain_sub,
		a130_send_history.from_service,
		a130_send_history.IP,
		a130_send_history.to_email,
		from_unixtime(a130_send_history.time_create,'\%Y-\%m-\%d \%H:\%i:\%s') AS sendtime
	FROM
		a130_send_history LEFT JOIN
		a130_received ON
		a130_received.body LIKE concat('".'%'."VERIFICATION-HASH: ',a130_send_history.hash,'%')
	WHERE
		a130_send_history.hash!='' AND
		a130_send_history.active='N' AND
		a130_received.active IS NULL
	";
	my $db0 = $main::DB{main}->Query( $sql );

	while ( my %db0_line = $db0->fetchhash )
	{
		# Chybajuci email

		# zistim prislusne emailove adresy, na ktore to treba poslat
		my $positions = 'manager;WEB';
		my $conf = `tom3-config_team --domain=$db0_line{domain}`;
		$db0_line{email} = "";
		foreach my $val ( split "\n", $conf )
		{
			next unless $val =~ /^\[ (\w+)\s+\] \[ \[ ([\w\@;\.-]+)\s+\]/;

			my $position = $1; my $email = $2;
			if ( $positions =~ /(^|;)$position(;|$)/ )
			{
				$db0_line{email} .= ';' if $db0_line{email};
				$db0_line{email} .= $email;
			}
		}
		my $db_send_to = $db0_line{email}; # emailove adresy - tvar pre databazu
		# tvar pre telo emailu
		$db0_line{email} = "<$db0_line{email}>"; $db0_line{email} =~ s|;|>,<|g;

		# odstranim email - forms@webcom.sk
		$db0_line{to_email} =~ s|forms\@webcom.sk||;
		$db0_line{to_email} =~ s/;+/;/g;
		$db0_line{to_email} =~ s/(^;|;$)//g;

		# naplnim info do designu
		while ( (my $k, my $v) = each %db0_line )
		{
			$XSGN{mail} =~ s|<%$k%>|$v|g;
		}
		$XSGN{mail} =~ s|<%.*?%>||g; # vycistim bordel

		# Poslem o tom email
		my $sql_send = "
		INSERT INTO TOM.a130_send
		SET
			sendtime='$main::time_current',
			priority='5',
			from_name='DEV',
			from_email='dev\@webcom.sk',
			from_host='server1.webcom.sk',
			from_service='a130 email checking',
			to_name='',
			to_email='$db_send_to',
			body='$XSGN{mail}'
		";
		#main::_log( $sql_send );
		$main::DB{main}->Query( $sql_send );
	}
	
}

1;

