# _mdl/130/130-control.forms-test-crash.cron
#!/bin/perl
# ÁÉÍÓÚ - USE UTF-8 !!!
package CRON::module;
use Net::HTTP::LiteAgent;
use Net::SMTP;
use CVML;
use strict;
use Mysql;
use Digest::MD5  qw(md5 md5_hex md5_base64);
use MIME::Base64;

our $authors = "gregor\@webcom.sk";

# tento cron kontroluje emaily ulozene v databaze, pochadzajuce z adresy forms@webcom.sk
sub execute
{
	my %env=@_;
	
	$env{err_email} = $TOM::contact_admin.";".$TOM::contact_operator.";".$TOM::contact_director.";".$TOM::contact_stats;
	$env{new_email} = $TOM::contact_admin.";".$TOM::contact_operator.";".$TOM::contact_director.";".$TOM::contact_stats;

	# spracovanie emailovych adries
	while ( (my $key, my $val) = each %env )
	{
		next unless $key =~ /^(err|new)_email$/;

		my $tmp_val;
		foreach my $email_addr ( sort split ";", $val )
		{
			$tmp_val .= "," if $tmp_val;
			$tmp_val .= "<$email_addr>";
		}
		$env{$key} = $tmp_val;
	}
	
	# tento cront je striktne 7day
	
	# CRON::debug::log(5, "Wow it runs");
	my $db0 = $main::DBH->Query("
	SELECT
		ID,
		name,
		domain,
		domain_sub,
		time_last,
		time_change
	FROM
		TOM.a130_services
	ORDER by domain ASC, domain_sub ASC
	");
	
	# time_last -> cas poslednej kontroly = spracovania udajov daneho riadku
	# time_change -> cas posledneho emailu  = zmeny udajov daneho riadku
	
	my %services;
	my $time_change = 0;
	while (my %db0_line = $db0->fetchhash)
	{
		my $path = "$db0_line{domain}/$db0_line{domain_sub}/$db0_line{name}";
		
		$services{$path}{time_last} = $db0_line{time_last};
		$services{$path}{count} = 0;
		$services{$path}{ID} = $db0_line{ID};
		
		# cas kedy som naposledy kontroloval services
		$time_change = $db0_line{time_last} unless $time_change;
		
	}
	
	my $db0 = $main::DBH->Query("
	SELECT
		rectime,
		body
	FROM
		TOM.a130_received
	WHERE
		rectime>$time_change
	ORDER BY rectime DESC");
	
	# premenna ktora mi povie o case posledneho emailu co som kontroloval
	# pokial taky nemam, pouzijem aktualny-600
	my $newtime;
	# prejdem si databazove emaily
	while (my %db0_line = $db0->fetchhash)
	{
		$newtime = $db0_line{rectime} unless $newtime;
		if ($db0_line{body} =~ /X-SERVICE: (.*?)\n/)
		{
			my $service = $1;
			$db0_line{body} =~ /X-DOMAIN: (.*?)\n/;
			my $domain = $1;
			$db0_line{body} =~ /X-DOMAIN-SUB: (.*?)\n/;
			my $domain_sub = $1;
			my $time_change = $db0_line{rectime};
			my $path = "$domain/$domain_sub/$service";
			
			# zapocitam mail a nastavim cas
			$services{$path}{count}++;
			$services{$path}{time_change} = $time_change;
		}
	}
	
	
	# nastavenie noveho casu
	$newtime = $cron::time_current-600 unless $newtime;
	
	
	# prejdem si spracovane info a updatnem databazu, pripadne rozoslem chybove emaily
	$env{err_message}=<<"HEADER";
From: a130-forms <dev\@webcom.sk>
To: <%email%>
Subject: Zrejme nefunkcny formular <%service%> na <%domain_sub%>

Toto je automaticky generovany email aplikacie na fontrolu formularov.
Z formularu <%service%> v domene <%domain_sub%> uz 7 dni neprisiel ziadny email.
HEADER

	$env{new_message}=<<"HEADER";
From: a130-forms <dev\@webcom.sk>
To: <%email%>
Subject: Novy formular <%service%> na <%domain_sub%>

Toto je automaticky generovany email aplikacie na fontrolu formularov.
Zaregistroval som novy formular <%service%> z domeny <%domain_sub%>
HEADER
	
	my @err_services;
	foreach my $serv_path (keys %services)
	{
		#$serv_path =~ /([^\/]+)\/([^\/]+)\/(.+)/;
		$serv_path =~ /^(.*)\/(.*)\/(.*)$/;
		
		my $service = $3; my $domain = $1; my $domain_sub = $2;
		
		if (!$services{$serv_path}{count})
		{
			main::_log(6, "Service $serv_path - possible damaged");
			main::_log(7, "Sending email to $env{err_email}");
			
			$env{err_message} =~ s|<%email%>|$env{err_email}|g;
			$env{err_message} =~ s|<%service%>|$service|g;
			$env{err_message} =~ s|<%domain%>|$domain|g;
			$env{err_message} =~ s|<%domain_sub%>|$domain_sub|g;
			
			$main::DBH->Query("
			INSERT INTO TOM.a130_send
			(
				sendtime,
				priority,
				from_name,
				from_email,
				from_host,
				from_service,
				to_name,
				to_email,
				body)
			VALUES	(
				'$cron::time_current',
				'5',
				'CRON',
				'info\@webcom.sk',
				'server1.webcom.sk',
				'a130-forms',
				'',
				'$env{err_email}',
				'$env{err_message}'
			)");
		}
		
		if (!$services{$serv_path}{ID})
		{
			# nova sluzba
			main::_log(6, "Got new service $service - $domain - $domain_sub");
			main::_log(7, "Sending email to $env{new_email}");
			
			$env{new_message} =~ s|<%email%>|$env{new_email}|g;
			$env{new_message} =~ s|<%service%>|$service|g;
			$env{new_message} =~ s|<%domain%>|$domain|g;
			$env{new_message} =~ s|<%domain_sub%>|$domain_sub|g;
			
			$main::DBH->Query("
			INSERT INTO TOM.a130_send
			(
				sendtime,
				priority,
				from_name,
				from_email,
				from_host,
				from_service,
				to_name,
				to_email,
				body)
			VALUES	(
				'$cron::time_current',
				'5',
				'CRON',
				'info\@webcom.sk',
				'server1.webcom.sk',
				'a130-forms',
				'',
				'$env{new_email}',
				'$env{new_message}'
			)");
			
			$main::DBH->Query("
				INSERT INTO TOM.a130_services
				(
					name,
					domain,
					domain_sub,
					email_count,
					time_last,
					time_create,
					time_change,
					active,
					lng
				)
				values
				(
					'$service',
					'$domain',
					'$domain_sub',
					$services{$serv_path}{count},
					$services{$serv_path}{time_last},
					$cron::time_current,
					$newtime,
					'Y',
					''
				)
			");
		}
		else
		{
			# sluzbu uz poznam
			$main::DBH->Query("
				UPDATE TOM.a130_services
				SET
					email_count=email_count+$services{$serv_path}{count},
					time_last=$newtime,
					time_change=$time_change
				WHERE
					ID=$services{$serv_path}{ID}
				LIMIT 1
			");
		}
	}
	
	return 1}

1;

